# branches/socket/scripts/zefirod.crontab.in
SHELL=/bin/bash
HOME=/
MAILTO=""

*/5 * * * 1-6 root /usr/local/sbin/zefirod.check
